# daemon/desklightd.service
[Unit]
Description=Desklight daemon

[Service]
Type=simple
ExecStart=/usr/bin/desklightd

[Install]
WantedBy=multi-user.target
